# shell/install/systemd/neuralvd.service
[Unit]
Description=NeuralV resident Linux daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/neuralvd
Restart=always
RestartSec=5
Environment=NEURALV_BASE_URL=https://neuralvv.org/basedata

[Install]
WantedBy=default.target
